# conf/systemd/controller.service.in
[Unit]
Description=Cloudbus Controller
After=network.target

[Service]
ExecStart=@prefix@/bin/controller
Type=exec

[Install]
WantedBy = multi-user.target
