# benchmarks/benchmark-four/conf/segment.service
[Unit]
Description=Cloudbus Segment
After=network.target

[Service]
ExecStart=/usr/local/bin/segment
Type=exec

[Install]
WantedBy = multi-user.target
